# ansible/roles/app/templates/systemd.socket.j2
[Unit]
Description={{ l3_domain }} socket

[Socket]
ListenStream={{ socket_path }}{{ socket_file }}

[Install]
WantedBy=sockets.target
